# worker/anpr-mqtt-worker.service
[Unit]
Description=ANPR MQTT Worker — auto-trigger /come on plate detection (Python)
Documentation=file:/opt/anpr-worker/README.md
# AlmaLinux service names: httpd (not apache2). The worker speaks only MQTT +
# backend HTTP — it never touches the DB directly — so the broker and the web
# server are its only real upstreams (no postgresql/mariadb dependency).
After=network-online.target mosquitto.service httpd.service
Wants=network-online.target mosquitto.service
# Bound the auto-restart so a hard failure can't hammer the box forever:
# if it dies > 5 times in 60s, stop trying and surface the failure.
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=anpr
Group=anpr
WorkingDirectory=/opt/anpr-worker
EnvironmentFile=-/opt/anpr-worker/.env
# Use the virtualenv created during deploy (see README)
ExecStart=/opt/anpr-worker/.venv/bin/python /opt/anpr-worker/worker.py
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=anpr-mqtt-worker

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/anpr-worker

[Install]
WantedBy=multi-user.target
